# src/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       

[supervisord]
logfile=/var/log/supervisord/supervisord.log    
logfile_maxbytes=50MB                           
logfile_backups=10                              
loglevel=error                                 
pidfile=/var/run/supervisord.pid               
nodaemon=false                                  
minfds=1024                                    
minprocs=200                                   
user=root                                      
childlogdir=/var/log/supervisord/


[program:linkedin_share]
command=python linkedin_share.python 
process_name=linkedin_share
autostart=true
autorestart=true
stdout_logfile=/var/log/linkedin_share/app.log
stderr_logfile=/var/log/linkedin_share/err.log

